# systemd/services/nas-update-weekly.service
[Unit]
Description=NAS Weekly Update & Maintenance 🔄
After=network.target local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nas_stage_update.sh
StandardOutput=journal
StandardError=journal
