# systemd/timers/nas-dlna-maint.timer
[Unit]
Description=NAS DLNA Maintenance (Daily 05:00 except Tuesday) 📺

[Timer]
# Run at 05:00 on all days EXCEPT Tuesday
OnCalendar=Mon,Wed,Thu,Fri,Sat,Sun *-*-* 05:00:00
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
