# systemd/timers/nas-borg-integrity.timer
[Unit]
Description=NAS Weekly Borg Integrity Check (Saturday 01:00) 🧬

[Timer]
# Run Saturdays at 01:00
OnCalendar=Sat *-*-* 01:00:00
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
